# .devcontainer/supervisord.conf
; supervisord.conf — Manages Xvfb, x11vnc, and websockify (noVNC)
;
; Started by start-emulator.sh. The emulator itself is launched by the script
; after supervisord brings up the display stack.

[supervisord]
nodaemon=false
logfile=/tmp/supervisord-emulator.log
pidfile=/tmp/supervisord-emulator.pid
childlogdir=/tmp

[program:xvfb]
command=Xvfb :99 -screen 0 1080x1920x24 -ac +extension GLX +render -noreset
autorestart=true
priority=10
stdout_logfile=/tmp/xvfb.log
stderr_logfile=/tmp/xvfb-error.log

[program:x11vnc]
command=x11vnc -display :99 -forever -shared -passwd minilauncher -rfbport 5900 -no6
autorestart=true
priority=20
stdout_logfile=/tmp/x11vnc.log
stderr_logfile=/tmp/x11vnc-error.log
depends_on=xvfb

[program:websockify]
command=websockify --web=/usr/share/novnc 6080 localhost:5900
autorestart=true
priority=30
stdout_logfile=/tmp/websockify.log
stderr_logfile=/tmp/websockify-error.log
depends_on=x11vnc